# cron/1-share-shadowsocks.cron
35 16 * * * /.../bin/python /.../parrotsec_bot/spider/share-sjadowsocks.py >> /.../parrotsec_bot/spider/logs.log 2>&1
