# cron/2-free-ss.cron
35 02 * * * /.../bin/python /.../free-ss_spider.py >> /.../logs.log 2>&1
