# spider/spider.cron
*/20 * * * * /home/.../bin/python /home/.../ParrotSecCN_Community_QQbot/spider/free-ss_spider.py >> /home/.../ParrotSecCN_Community_QQbot/spider/logs.log 2>&1
35 16 * * * /home/.../bin/python /home/.../ParrotSecCN_Community_QQbot/spider/ss_pythonic_spider.py >> /home/.../ParrotSecCN_Community_QQbot/spider/logs.log 2>&1
